# corobot_bringup/corobot-kinect.conf
description	"Bringup script for the corobot kinect nodes."

start on started corobot
stop on stopping corobot

console output
respawn

exec corobot-start-kinect
